# omnibus/package-scripts/gemfast/postinst
#!/bin/sh

useradd -M -U gemfast
groupadd -f gemfast
usermod -a -G gemfast gemfast

is_darwin()
{
  uname -a | grep "^Darwin" 2>&1 >/dev/null
}

if is_darwin; then
    PREFIX="/usr/local"
    mkdir -p "$PREFIX/bin"
else
    PREFIX="/usr"
fi

INSTALLER_DIR=/opt/gemfast

rm -f $PREFIX/bin/gemfast

ln -sf $INSTALLER_DIR/bin/gemfast $PREFIX/bin || error_exit "Cannot link gemfast to $PREFIX/bin"

mkdir -p /var/gemfast/db
chown gemfast:gemfast /var/gemfast
chown gemfast:gemfast /etc/gemfast

for sv in "caddy" "gemfast"
do
  cp "/opt/gemfast/systemd/$sv/$sv.service" /etc/systemd/system
  systemctl enable $sv --now
done

echo "Thank you for installing Gemfast Server! For help getting started visit https://gemfast.io/docs"

exit 0